# ansible/backup/files/dbbackup.service
[Unit]
Description=Daily DB Table Backup Service
After=network.target

[Service]
User=dbbackup
Group=dbbackup
ExecStart=serpent-surge-backup
Nice=10

[Install]
WantedBy=multi-user.target
